# ansible/files/archives/opt-custom/smf/rc-pre-network.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='system:pre-network'>
<service
        name='system/postboot-rc/pre-network'
        type='service'
        version='1'>

        <create_default_instance enabled='true' />

        <single_instance />

        <dependency
                name='fs-smartdc'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/system/filesystem/smartdc' />
        </dependency>

        <dependent
                name='net-lo'
                grouping='require_all'
                restart_on='none' >
                <service_fmri value='svc:/network/loopback' />
        </dependent>

        <dependent
                name='net-physical'
                grouping='require_all'
                restart_on='none' >
                <service_fmri value='svc:/network/physical' />
        </dependent>

        <dependent
                name='net-ipf'
                grouping='require_all'
                restart_on='none' >
                <service_fmri value='svc:/network/ipfilter' />
        </dependent>

        <dependent
                name='net-ipsec'
                grouping='require_all'
                restart_on='none' >
                <service_fmri value='svc:/network/ipsec/ipsecalgs' />
        </dependent>

        <exec_method
                type='method'
                name='start'
                exec='/opt/custom/method/postboot-rc-init start pre-network'
                timeout_seconds='300'>
        </exec_method>

        <exec_method
                type='method'
                name='stop'
                exec='/opt/custom/method/postboot-rc-init stop pre-network'
                timeout_seconds='300'>
        </exec_method>

        <exec_method
                type='method'
                name='refresh'
                exec='/opt/custom/method/postboot-rc-init refresh pre-network'
                timeout_seconds='300'>
        </exec_method>

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
                <propval name='critical_failure_count' type='integer' value='1'/>
        </property_group>

        <stability value='Unstable' />

</service>
</service_bundle>
